# mimo-free-proxy.service.example
[Unit]
Description=MiMo Free OpenAI-Compatible Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mimo-free-proxy
EnvironmentFile=/etc/mimo-free-proxy.env
ExecStart=/opt/mimo-free-proxy/mimo-free-proxy
Restart=on-failure
RestartSec=3
KillMode=control-group
TimeoutStopSec=15
User=root

[Install]
WantedBy=multi-user.target
